# files/jenkins-startup.service
[Unit]
Description=Jenkins Startup Configuration
After=network-online.target jenkins.service
Wants=network-online.target
Requires=jenkins.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure-caddy.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target